# check-all-vpn-exits.cron
*/10 * * * * root /usr/lib/gatemon/check-all-vpn-exits.sh || echo "check exited with error"
