# check-all-vpn-exits.cron
*/10 * * * * root /opt/gatemon/check-all-vpn-exits.sh
