# check-all-vpn-exits.cron
*/10 * * * * root /usr/bin/check-all-vpn-exits.sh || echo "check exited with error"
