# gatemon.cron
*/10 * * * * root /opt/gatemon/gatemon.sh
