# gatemon.root.service
[Unit]
Description=gatemon (Check all services of every gateway)
After=network-online.target

[Service]
Type=simple
ExecStart=/opt/gatemon/gatemon.sh
TimeoutStopSec=300
KillMode=process
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
